# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

# ensure the /etc/cron.d dir exists
mkdir -p /etc/cron.d

# -f	Foreground
#-d N	Set log level, log to stderr. Most verbose 0, default 8
# -c DIR	Cron dir

crond -f -d 6
